# files/etc/uci-defaults/99-init-settings
#!/bin/bash
#=================================================
# File name: init-settings.sh
# Description: This script will be executed during the first boot
# Author: SuLingGG
# Blog: https://mlapp.cn
#=================================================

# Set default theme to luci-theme-argon
uci set luci.main.mediaurlbase='/luci-static/argon'

# 设置DHCP服务器
uci set dhcp.@dnsmasq[0].local='/ycslan/'
uci set dhcp.@dnsmasq[0].domain='ycslan'

# 其他自定义设置
# 禁用DHCP服务
uci set dhcp.lan.ignore='1'
# 禁用DHCP_IPv6服务
uci set dhcp.lan.dhcpv6='disabled'
uci set dhcp.lan.ra='disabled'
uci set dhcp.lan.ndp='disabled'
# 下面的为示例命令，你可使用。如若修改，更改变量值中的端口号即可，务必不能出现端口冲突！80和443端口将交与Lucky监听，因此必须修改。
uci set uhttpd.main.listen_http='0.0.0.0:72'
uci set uhttpd.main.listen_https='0.0.0.0:445'
