# files/etc/uci-defaults/99-init-settings
#!/bin/bash
#=================================================
# File name: init-settings.sh
# Description: This script will be executed during the first boot
# Author: SuLingGG
# Blog: https://mlapp.cn
#=================================================

# Set default theme to luci-theme-argon
uci set luci.main.mediaurlbase='/luci-static/argon'

# 设置DHCP服务器
uci set dhcp.@dnsmasq[0].local='/ycslan/'
uci set dhcp.@dnsmasq[0].domain='ycslan'

# 其他自定义设置
# 禁用DHCP服务
uci set dhcp.lan.ignore='1'
# 禁用DHCP_IPv6服务
uci set dhcp.lan.dhcpv6='disabled'
uci set dhcp.lan.ra='disabled'
uci set dhcp.lan.ndp='disabled'

# 旁路由SNAT(这个必须加上)
uci add firewall nat
uci set firewall.@nat[-1].name='SideRouter'
uci add_list firewall.@nat[-1].proto='all'
uci set firewall.@nat[-1].src='lan'
uci set firewall.@nat[-1].target='SNAT'
uci set firewall.@nat[-1].snat_ip='200.56.72.248' # 这里填上旁路由IP
exit 0
